# src/package/funjsq/files/etc/init.d/funjsq
#!/bin/sh /etc/rc.common

START=62

boot() {
	module_name=`cat module_name`

	[ "x$module_name" == "xR7800" ] && {
		/data/funjsq/bin/funjsq.sh init &
	}

}
