# src/package/base-files/files/etc/rc.d/S62funjsq
#!/bin/sh /etc/rc.common


START=62

boot() { 
	/data/funjsq/bin/funjsq.sh init

	/data/funjsq/bin/FYdaemon
}
